# deploy/system-conf/etc/profile.d/local.sh
export PATH="$HOME/.local/bin:$PATH"
export PATH="$PATH:/miniconda/bin"
export LD_LIBRARY_PATH="${HOME}/.local/lib:${LD_LIBRARY_PATH}"

source activate py2
# eval "$(pyenv init -)"
